# iso-profiles/netrunner/kde/live-overlay/usr/lib/systemd/system/autologin-into-wayland.service
[Unit]
Description=LiveMedia Autologin into Wayland Session
Before=systemd-user-sessions.service

[Service]
Type=oneshot
ExecStart=/usr/bin/autologin-into-wayland

[Install]
WantedBy=multi-user.target
